# ops/systemd/certistock-backup-weekly.timer
[Unit]
Description=Run CertiStock weekly backup

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
