# ops/systemd/certistock-backup-daily.service
[Unit]
Description=CertiStock daily backup

[Service]
Type=oneshot
User=postgres
EnvironmentFile=-/opt/apps/certistock/shared/env
ExecStart=/usr/local/bin/certistock-backup.sh daily
