# ops/systemd/certistock-retention.timer
[Unit]
Description=Run CertiStock retention automation

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
